# ctf.service
[Unit]
Description=cyber262-ctf

[Service]
ExecStart=/root/cyber262-ctf/run.sh
ExecStop=/bin/kill -9 `cat /root/cyber262-ctf/run.pid`
Type=simple
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
